# templates/unicorn-init.erb
#! /bin/sh

# File: /etc/init.d/unicorn

### BEGIN INIT INFO
# Provides:          unicorn
# Required-Start:    $local_fs $remote_fs $network $syslog
# Required-Stop:     $local_fs $remote_fs $network $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: starts the unicorn web server
# Description:       starts unicorn
### END INIT INFO

COMMAND_OPTS="-c <%= config_path %>/unicorn.rb -E production -D"
PID="<%= pid_path %>/unicorn.pid"
NAME=unicorn

APP_ROOT="<%= app_root %>"
COMMAND="bundle exec unicorn_rails"
cd $APP_ROOT || exit 1

case "$1" in
  start)
	echo -n "Starting $NAME: "
	su -c "cd $APP_ROOT && $COMMAND $COMMAND_OPTS" - deploy
	echo "[OK]"
	;;
  stop)
	echo -n "Stopping $NAME: "
        kill -QUIT `cat $PID`
	echo "[OK]"
	;;
  restart)
	echo -n "Restarting $NAME: "
        kill -USR2 `cat $PID`
	echo "[OK]"
	;;
  reload)
        echo -n "Reloading $NAME configuration: "
        kill -USR1 `cat $PID`
        echo "[OK]"
        ;;
  *)
	echo "Usage: $NAME {start|stop|restart|reload}" >&2
	exit 1
	;;
esac

exit 0
